# install-trading/novnc/fz-grid/systemd/fz-grid@.service
[Unit]
Description=FZ-Grid Playwright Runner für Instanz %i
After=fz-xvfb@%i.service
Requires=fz-xvfb@%i.service

[Service]
MemoryHigh=1800M
EnvironmentFile=/opt/fz-grid/env/%i.env
WorkingDirectory=/opt/fz-grid
ExecStart=/usr/bin/node /opt/fz-grid/runner.js
ExecStopPost=/opt/scripts/session-end.sh %i
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
