# install-trading/novnc/fz-grid/systemd/fz-x11vnc@.service
[Unit]
Description=x11vnc für FZ-Grid Instanz %i
After=fz-xvfb@%i.service
Requires=fz-xvfb@%i.service

[Service]
EnvironmentFile=/opt/fz-grid/env/%i.env
ExecStart=/usr/bin/x11vnc -display ${DISPLAY} -rfbport ${VNC_PORT} -nopw -forever -shared -quiet
Restart=on-failure
RestartSec=3
SuccessExitStatus=2 15

[Install]
WantedBy=multi-user.target
