# install-trading/novnc/maintenance/fz-grid-shutdown-check.service
[Unit]
Description=FZ-Grid Shutdown Check
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/opt/maintenance/fz-grid-shutdown-check.sh
TimeoutStartSec=300
