# install-trading/guacamole/maintenance/update-and-shutdown.service
[Unit]
Description=System-Update und Shutdown
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/maintenance/update-and-shutdown.sh
TimeoutStartSec=1800
